# scripts/systemd/ombrouter-apikey.service.example
# Copy to ~/.config/systemd/user/ombrouter-apikey.service
# Chat proxy: Bearer to your OpenAI-compatible upstream.
# 1. Create ~/.config/ombrouter/apikey.env (mode 600) with:
#    OMBROUTER_UPSTREAM_API_BASE=https://api.openai.com/v1
#    OMBROUTER_UPSTREAM_API_KEY=sk-...
# 2. Edit ExecStart if ombrouter is not on PATH.

[Unit]
Description=OmbRouter apiKey proxy (BYOK chat)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=BLOCKRUN_PROXY_PORT=8403
Environment=OMBROUTER_UPSTREAM_MODE=apiKey
EnvironmentFile=-%h/.config/ombrouter/apikey.env
Restart=on-failure
RestartSec=5
ExecStart=ombrouter

[Install]
WantedBy=default.target
